# frontend/device/orangepi-capture.service
[Unit]
Description=Orange Pi Capture & Roboflow Service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=simple
# Adjust paths below
WorkingDirectory=/opt/cereal-device
EnvironmentFile=/opt/cereal-device/.env
ExecStart=/opt/cereal-device/venv/bin/python -u orangepi_capture.py
Restart=on-failure
RestartSec=5
# Optional: give access to video & serial groups
User=orangepi
Group=orangepi

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
ReadWriteDirectories=/opt/cereal-device

[Install]
WantedBy=multi-user.target
